# vlog-daily.timer
[Unit]
Description=Run VRChat Vlog Daily Processing at 6AM JST

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true
Unit=vlog-daily.service

[Install]
WantedBy=timers.target
